# code/开机自启与服务保活/elf-glasses-net-failover.service
[Unit]
Description=ELF Glasses WiFi preferred / 4G fallback route keeper
After=network.target

[Service]
Type=simple
ExecStart=/home/elf/elf_glasses_net_failover.sh run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
